# config/home/uroborosq/.config/autostart/preventsuspendgnome.desktop
[Desktop Entry] 
Type=Application
Exec="/home/uroborosq/Рабочий стол/Linux/cmd/bin/playerctl-watch"
