# config/usr/lib/systemd/system/cpuboostmanager.service
[Unit]
Description=Foo

[Service]
ExecStart="/home/uroborosq/Рабочий стол/Linux/cmd/bin/cpu-boost-manager"

[Install]
WantedBy=multi-user.target
