# config/usr/lib/systemd/system/cpuboostmanager.service
[Unit]
Description=Foo

[Service]
ExecStart="/usr/bin/cpu-boost-manager"

[Install]
WantedBy=multi-user.target
